# etc/cron.d/oss.jobs
* * * * * root /usr/share/oss/tools/check-salt-master.sh
