# setup/systemd/firewalld.service.d/FirewalldExecStart.conf
[Service]
Restart=always
RestartSec=10
ExecStartPre=/usr/share/cranix/tools/sync-rooms-to-firewalld.py
ExecStartPost=/usr/share/cranix/tools/set_fw_default_status.sh

